# files/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.
# Untuk Mengaktifkan Hapus Tanda Pagar (#) #
# Auto Sinkron Jam dengan Bug #
# Example Operator XL #
# /usr/bin/jam.sh bsns.bilibili.tv:80
# auto hapus chache 8 jam sekali #
# * */8 * * * /sbin/free.sh
exit 0
